# stunserver.conf
start on runlevel [2345]
stop on runlevel [016]

exec /home/hyun/github/stunserver/stunserver -v 2
